# katalogue/systemd/katalogued.service
[Unit]
Description=Katalogue daemon

[Service]
Type=simple
ExecStart=/usr/bin/katalogued

[Install]
WantedBy=default.target
